# ami/tsc-clocksource.conf
#!upstart

start on runlevel [2345]

task

script
  exec echo tsc > /sys/devices/system/clocksource/clocksource0/current_clocksource
end script

