# ami/xenial/docker-swapoff.crontab
# disable docker swap every minute
* * * * * root /usr/bin/flock -n /var/lock/docker-swapoff.lock /usr/local/bin/docker-swapoff
